# images/php56src_apache_sendmail_visor/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=service apache2 start

[program:sendmail]
command=service sendmail start
