# images/php56_apache_debian_visor/supervisord.conf
[supervisord]
nodaemon=true

#[program:startup]
#command=/opt/yurii-tools/startup.sh

[program:apache2]
command=service apache2 start

[program:sendmail]
command=service sendmail start

